# systemd/qrapp.service
[Unit]
Description=Gunicorn instance to serve QR App
After=network.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
# User/Group: Dedicated service user for production
User=qrapp
Group=www-data

# Directories managed by systemd
# StateDirectory -> /var/lib/qrapp (Persistent data)
# RuntimeDirectory -> /run/qrapp (Sockets, Locks)
# LogsDirectory -> /var/log/qrapp (Logs)
StateDirectory=qrapp
RuntimeDirectory=qrapp
RuntimeDirectoryMode=0755
LogsDirectory=qrapp

WorkingDirectory=/opt/qr_code_backend
Environment="PATH=/opt/qr_code_backend/.venv/bin:/usr/bin"
EnvironmentFile=/etc/qr-secrets.env

# Strict Production Paths
Environment=INSTANCE_DIR=/var/lib/qrapp
Environment=PRINT_INBOX_DIR=/var/lib/qrapp/print_inbox

# Run migrations exactly once before startup
ExecStartPre=/opt/qr_code_backend/.venv/bin/python3 migrate.py

# Safe gunicorn start
# --bind unix:/run/qrapp/qrapp.sock : Absolute path to RuntimeDirectory
# -m 007 : Safe umask for socket group access
ExecStart=/opt/qr_code_backend/.venv/bin/gunicorn --workers 3 --bind unix:/run/qrapp/qrapp.sock -m 007 app:app

Restart=on-failure
RestartSec=2
TimeoutStartSec=60
TimeoutStopSec=30

# =============================================================================
# Security Hardening (Correct Forever)
# =============================================================================
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/qrapp /run/qrapp
UMask=0007

[Install]
WantedBy=multi-user.target
